# crontab
0 0 * * * root cd /usr/src/app && python3 manage.py update_csv